# setup_files/viz_dashboard.service
[Unit]
Description=Runs a visualization dashboard for the collected sensor data.
After=network.target  mqtt.service

[Service]
ExecStart=/home/beta/sensor-hub/bin/run_sensor_dashboard
# ExecStart=/home/beta/sensor-hub/bin/gunicorn -w 4 -b 0.0.0.0:8050 sensor_dashboard.run_app:server
WorkingDirectory=/home/beta/sensor-hub/lib/python3.11/site-packages/sensor_dashboard
Restart=always

[Install]
WantedBy=multi-user.target

